# scripts/units/soletta-dev-app-server.service.in
[Unit]
Description=Run Soletta Development Application server
After=warn-boot.service

[Service]
ExecStart=/usr/bin/NODE_BIN_NAME PATH/server/app.js 1> /dev/null 2>&1
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target
